# runit/optimus-manager/finish
#!/bin/sh
exec /usr/bin/python3 -u -m optimus_maanger.hooks.post_daemon_stop

