# runit/optimus-manager/run
#!/bin/sh
sv o pre_start
exec /usr/bin/python3 -u -m optimus_manager.daemon
