# runit/prime-switch-boot/run
#!/bin/sh
sv o prime-switch-boot
exec /usr/bin/prime-switch-boot
