# systemd/optimus-manager.service
[Unit]
Description=Optimus Manager
Before=display-manager.service

[Service]
ExecStart=/usr/bin/python -u /usr/bin/optimus-manager-daemon --startup
Restart=always

[Install]
WantedBy=multi-user.target
