# systemd/optimus-manager.service
[Unit]
Description=Optimus Manager Setup Service
Before=display-manager.service
BindsTo=display-manager.service
Wants=optimus-manager-daemon.service

[Service]
Type=oneshot
ExecStart=/usr/bin/python -u /usr/bin/optimus-manager-setup --setup-start
ExecStop=/usr/bin/python -u /usr/bin/optimus-manager-setup --setup-stop
RemainAfterExit=yes

[Install]
WantedBy=graphical.target
