# runit/pre_start_optimus-manager/run
#!/bin/sh
exec /usr/bin/python3 -u -m optimus_manager.pre_daemon_start
