# package/optimus-manager.install
#! /bin/bash


pre_install () {
	rm --force \
		"etc/X11/xorg.conf" \
		"etc/X11/nvidia-xorg.conf.d/"*nvidia*
}


post_install () {
	systemctl enable "optimus-manager.service"
}


post_upgrade () {
	systemctl daemon-reload
}


post_remove() {
	systemctl disable --now "optimus-manager.service"
	rm --force "etc/X11/xorg.conf.d/10-optimus-manager.conf"
}
