# package/optimus-manager.install
#! /bin/bash


post_install() {
	systemctl enable "optimus-manager.service"
}


post_upgrade () {
	systemctl daemon-reload
}


post_remove() {
	systemctl disable "optimus-manager.service"
	rm --force "/etc/X11/xorg.conf.d/10-optimus-manager.conf"
}
